# recipes-apps/my-cpp-app/files/my-cpp-app.service
[Unit]
Description=My C++ Application
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/my-cpp-app
Restart=always
RestartSec=5
WatchdogSec=30
NotifyAccess=all

[Install]
WantedBy=multi-user.target

